# aircast2/rootfs/etc/s6-overlay/s6-rc.d/nqptp/finish
#!/command/with-contenv bashio
# ==============================================================================
# Handle NQPTP exit
# ==============================================================================
declare exit_code=${1}
declare signal=${2}

bashio::log.info "NQPTP exited (code=${exit_code}, signal=${signal})"

if [[ "${exit_code}" -ne 0 ]] && [[ "${exit_code}" -ne 256 ]]; then
    bashio::log.warning "NQPTP crashed — restarting..."
fi
