# aircast2/rootfs/etc/s6-overlay/s6-rc.d/nqptp/run
#!/command/with-contenv bashio
# ==============================================================================
# Start NQPTP - PTP timing daemon required for AirPlay 2
# ==============================================================================
bashio::log.info "Starting NQPTP timing daemon..."
exec nqptp
