# aircast2/rootfs/etc/s6-overlay/s6-rc.d/mosquitto/finish
#!/command/with-contenv bashio
# ==============================================================================
# Handle Mosquitto exit
# ==============================================================================
declare exit_code=${1}
declare signal=${2}

bashio::log.info "Mosquitto exited (code=${exit_code}, signal=${signal})"

if [[ "${exit_code}" -ne 0 ]] && [[ "${exit_code}" -ne 256 ]]; then
    bashio::log.warning "Mosquitto crashed — restarting..."
fi
