# aircast2/rootfs/etc/s6-overlay/s6-rc.d/bridge/finish
#!/command/with-contenv bashio
# ==============================================================================
# Handle bridge exit
# ==============================================================================
declare exit_code=${1}
declare signal=${2}

bashio::log.info "Bridge exited (code=${exit_code}, signal=${signal})"

if [[ "${exit_code}" -ne 0 ]] && [[ "${exit_code}" -ne 256 ]]; then
    bashio::log.error "Bridge crashed — check logs above for Python traceback"
fi
